# server/config/supervisord.conf
[supervisord]
nodaemon=true

[program:init_db]
command=bash /opt/init_db.sh

[program:tomcat]
command=bash /opt/tomcat/bin/catalina.sh start
